# followsun.service
[Unit]
Description=FollowSun - GNOME theme switcher based on sunrise/sunset
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment="CONFIG_DIR=%h/.config/followsun"
ExecStartPre=/bin/sh -c 'mkdir -p ${CONFIG_DIR}'
ExecStart=/home/tomas/dev/bash/followsun/followsun.sh --daemon
Restart=on-failure
RestartSec=30

[Install]
WantedBy=graphical-session.target