# followsun.timer
[Unit]
Description=Run FollowSun every 5 minutes

[Timer]
OnBootSec=30sec
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
